# src/theme-switcher/io.starikov.theme-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.starikov.theme-watcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>bash</string>
        <string>/Users/starikov/.dotfiles/src/theme-switcher/auto-theme-watcher.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ProcessType</key>
    <string>Background</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/starikov</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/starikov/.cache/theme-switcher/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/starikov/.cache/theme-switcher/launchd.error.log</string>
</dict>
</plist>